# minecraft.cron
0 4 * * * root /path/to/backup/script/mc-backup.sh